# scripts/systemd/atlas-pool.service
[Unit]
Description=Atlas LVM thin pool
# The pool lives on a loopback PV (a sparse backing file on the root fs). The
# backing file survives a reboot, but the loop device binding does not, so the
# VG/pool are gone until the loop is re-attached. bootstrap-server.sh creates
# the pool but is NOT re-run on boot, so this oneshot re-asserts it. Ordered
# before any VM unit (firecracker-vm@.service Wants network-online.target;
# this only needs local fs) so a VM never starts against a missing pool.
After=local-fs.target
Before=firecracker-vm@.service
DefaultDependencies=no
RequiresMountsFor=/var/lib/atlas/pool

[Service]
Type=oneshot
RemainAfterExit=yes
# Re-assert the pool via the durable atlas package. ThinPool().ensure() is
# idempotent: if the pool LV already exists it only re-attaches the loop device
# and activates the VG; otherwise it rebuilds pv/vg/pool from the backing file.
# The package lives under /var/lib/atlas/bin (placed by bootstrap), added to
# sys.path here so `from atlas.lvm import ThinPool` resolves.
ExecStart=/usr/bin/python3 -c "import sys; sys.path.insert(0, '/var/lib/atlas/bin'); from atlas.lvm import ThinPool; ThinPool().ensure()"

[Install]
WantedBy=multi-user.target
