# proxy/guest/nginx.service.d/atlas.conf
# Atlas drop-in over the nginx.org package's own /lib/systemd/system/nginx.service
# (proxy-design.md §8). We do NOT shadow the whole unit — it already ships
# Type=forking, PIDFile, ExecStart=-c ${CONFFILE} (CONFFILE defaults to our
# /etc/nginx/nginx.conf), ExecReload (kill -HUP), ExecStop and
# WantedBy=multi-user.target, so `apt upgrade nginx` keeps shipping base-unit fixes.
# This carries ONLY the deltas with no stock equivalent.
[Unit]
# Order nginx after the guest's static /128 IPv6 is up, else the v6 listeners hit
# EADDRNOTAVAIL on a fresh boot.
After=atlas-network.service

[Service]
# The package unit ships NO config precheck. Re-add it: with Restart=on-failure
# below, a config that fails `nginx -t` must refuse to START (clear preflight fail)
# rather than die and restart-loop. (The leading empty ExecStartPre= clears any
# value a base-unit upgrade might add, so ours is the only one.)
ExecStartPre=
ExecStartPre=/usr/sbin/nginx -t -c /etc/nginx/nginx.conf
# The ~20000-listener pool needs the fd ceiling above systemd's ~1024 default; the
# conf sets no worker_rlimit_nofile. The per-VM cgroup caps (vcpus, memory) are the
# real resource bound, sized on the proxy VM — not pinned here.
LimitNOFILE=1048576
Restart=on-failure
# Create /run/nginx (tmpfs, stock ships no such dir) at 0750 root before ExecStart,
# so the admin socket (listen unix:/run/nginx/admin.sock) lands in a root-only dir.
# This subsumes the standalone tmpfiles.d file. (The pidfile /run/nginx.pid sits
# loose in /run, which systemd always provides.)
RuntimeDirectory=nginx
RuntimeDirectoryMode=0750
